# chat-server/chat-server.service
[Unit]
Description=Chat Server Service
After=network.target

[Service]
Type=simple
User=user
ExecStart=/home/user/orbyt/chat-server/start-server.sh
Restart=on-failure
RestartSec=10
StandardOutput=append:/var/log/chat-server.log
StandardError=append:/var/log/chat-server.log

[Install]
WantedBy=multi-user.target
